# .config/systemd/user/astal-watcher.path
[Path]
PathModified=%h/.config/astal
PathModified=%h/.config/astal/widget

[Install]
WantedBy=hyprland-session.target
