# .config/systemd/user/astal-watcher.path
[Path]
PathModified=%h/.config/astal

[Install]
WantedBy=hyprland-session.target
