# .config/systemd/user/update_lsps.service
[Unit]
Description=Download newest versions of globally installed LSPs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dotfiles/
ExecStart=make pnpm
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
