# .config/mako/mako-watcher.path

# YES, you need to enable this path service as well !

[Path]
PathModified=/home/earthian/.config/mako

[Install]
WantedBy=hyprland-session.target
