# .config/swaync/swaync-watcher.path

# YES, you need to enable this path service as well !

[Path]
PathModified=/home/earthian/.config/swaync

[Install]
WantedBy=hyprland-session.target
